# vault.service
[Unit]
Description=Vault dev server
          
[Service]
Type=simple
ExecStart=/home/${user}/vault server -dev -dev-listen-address=0.0.0.0:${port} 